# rsync-ssh/services/ssh/run
#!/bin/sh

autossh -M 22874 -o StrictHostKeyChecking=no -i /local/id_rsa -4N -R 22873:127.0.0.1:873 $REMOTE
