# rsync-ssh/services/rsyncd/run
#!/bin/sh

/usr/bin/rsync --no-detach --daemon --log-file=/dev/stdout --config=/local/rsyncd.conf